# supervisord-mysqld.conf
[program:mysqld]
#command=/start-mysqld.sh
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --pid-file=/var/run/mysqld/mysqld.pid
numprocs=1
autostart=true
autorestart=true
